# Curl/crontest.cron
*/1 * * * * /bin/date  >> /Users/george/jdlog.txt