# Curl/crontest.cron
*/6 * * * *  /Users/george/Code/HappyLearning/Curl/jd.sh  >> ~/jd_log.txt